# bin/lambdacd.conf
description "Run a LambdaCD pipeline"

start on runlevel startup
stop on runlevel shutdown

respawn

exec java -jar /opt/pipeline/pipeline-0-standalone.jar