# resources/deployscripts/demo-app.conf
description "Run my project's jar"

start on runlevel startup
stop on runlevel shutdown

respawn

exec java -jar /opt/app/app.jar